# site-cookbooks/site-thinking-sphinx/templates/default/thinking_sphinx.conf.erb
description "Sphinx Upstart Script"

start on startup
stop on shutdown

console output

respawn

script
    env RAILS_ENV=production
    exec sudo -i /usr/local/bin/searchd --config <%= node[:site_app][:root] %>/config/production.sphinx.conf --nodetach
end script
