# site-cookbooks/site-app/templates/default/thinking_sphinx.conf.erb
description "Thinking sphinx worker configuration. Run with ID"

start on (local-filesystems and net-device-up IFACE=eth0)
stop on shutdown

respawn
respawn limit 5 20

script
  PIDFILE=/var/run/thinking_sphinx.pid
  echo $$ > $PIDFILE
  chown www-data:www-data $PIDFILE
  exec su -c "cd <%= node[:site_app][:root] %> && /opt/rbenv/shims/bundle exec rake ts:start RAILS_ENV=production && /opt/rbenv/shims/bundle exec rake ts:rebuild RAILS_ENV=production >> /var/log/thinking_sphinx.log 2>&1" www-data
end script